# OpenRC/chains.init
#!/sbin/openrc-run

name="chains"
description="A transparent proxy v2ray/xray chains client by ru osint team - tw"
command="/usr/local/bin/chains"
pidfile="/run/${RC_SVCNAME}.pid"
command_background="yes"
rc_ulimit="-n 30000"
rc_cgroup_cleanup="yes"

 depend() {
	need net
	after net
   }
